# .config/systemd/user/trash-empty.timer
[Unit]
Description=Empty trash older than 30 days daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
